# scripts/github-runner/github-runner.service
[Unit]
Description=GitHub Actions Ephemeral Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=/etc/github-runner/runner.env
ExecStart=/usr/local/bin/github-runner-bootstrap
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/actions-runner /etc/github-runner
WorkingDirectory=/opt/actions-runner

[Install]
WantedBy=multi-user.target
